# deploy/systemd/trackerbundle-ebay-scheduler.service
[Unit]
Description=TrackerBundle eBay ISBN Scheduler
After=network-online.target trackerbundle-api.service
Wants=network-online.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/trackerbundle3
Environment="PATH=/home/ubuntu/trackerbundle3/venv/bin"
EnvironmentFile=/etc/trackerbundle.env
ExecStart=/home/ubuntu/trackerbundle3/venv/bin/python -m app.scheduler_ebay
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
